# rootfs/etc/services.d/neso-watch/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Start the NESO watch service
# ==============================================================================

bashio::log.info "Starting NESO Octopus Watch..."

# Verify environment
bashio::log.info "Checking configuration..."
if bashio::config.has_value 'scan_interval'; then
    export SCAN_INTERVAL=$(bashio::config 'scan_interval')
    bashio::log.info "Scan interval set to: ${SCAN_INTERVAL}"
else
    bashio::log.warning "No scan_interval set, using default: 300"
    export SCAN_INTERVAL=300
fi

export DEBUG=$(bashio::config 'debug')

# Verify Python installation
bashio::log.info "Checking Python installation..."
if ! command -v python3 &> /dev/null; then
    bashio::log.error "Python3 not found!"
    exit 1
fi

# Verify script exists
if [ ! -f "/usr/src/app/nesoscan.py" ]; then
    bashio::log.error "nesoscan.py not found in /usr/src/app!"
    ls -la /usr/src/app
    exit 1
fi

# Verify data directory
if [ ! -d "/data" ]; then
    bashio::log.warning "Data directory not found, creating..."
    mkdir -p /data
fi

cd /usr/src/app
bashio::log.info "Starting main script..."
exec python3 nesoscan.py